# xorcism_ts/connectors/detections-ai/sample.json
{
  "_comment": "Illustrative shape of the detections.ai BFF /bff/api/v1/intel payload. Field names are normalized tolerantly by run.py (_normalize_one). Used for offline import and `runner.py --selftest --connector detections-ai`.",
  "data": {
    "items": [
      {
        "id": "edrchoker-qos-throttling",
        "slug": "edrchoker-qos-throttling",
        "title": "EDRChoker Bypasses EDR Telemetry via QoS Throttling",
        "summary": "The EDRChoker tool weaponizes Windows Policy-based Quality of Service (QoS) to throttle EDR process bandwidth to 8 bits per second, starving sensors of telemetry while attacker activity proceeds.",
        "url": "https://detections.ai/intel-exchange/edrchoker-qos-throttling",
        "publishedAt": "2026-06-12T09:15:00Z",
        "author": { "name": "Tim Peck", "handle": "timpeck" },
        "views": 244,
        "techniques": [
          { "id": "T1562.001", "name": "Impair Defenses: Disable or Modify Tools" },
          { "id": "T1071.001", "name": "Application Layer Protocol: Web Protocols" }
        ],
        "actors": ["TeamPCP"],
        "malware": ["EDRChoker"],
        "cves": [],
        "tags": ["edr-evasion", "windows", "defense-evasion"]
      },
      {
        "id": "dprk-ironchain-backdoor",
        "slug": "dprk-ironchain-backdoor",
        "title": "DPRK-linked Ironchain Backdoor Targets Supply Chains",
        "description": "A Lazarus-linked cluster deploys the Ironchain backdoor through a poisoned build dependency, exploiting CVE-2026-41089 for initial access and using scheduled tasks for persistence.",
        "slugUrl": null,
        "createdAt": "2026-06-09T22:40:00Z",
        "author": "threatlabs",
        "views": 263,
        "mitreTechniques": ["T1053.005", "T1195.001", "T1059.001"],
        "adversaries": ["Lazarus Group"],
        "tools": ["Ironchain", "Cobalt Strike"],
        "vulnerabilities": ["CVE-2026-41089"],
        "tags": ["dprk", "supply-chain", "persistence"]
      }
    ]
  }
}
